# packaging/systemd/uroam.service
[Unit]
Description=UROAM - Unified RAM Optimization and Management Framework
Documentation=man:uroamd(8)
After=multi-user.target
Requires=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/uroamd
Restart=on-failure
RestartSec=5
RuntimeDirectory=uroam
RuntimeDirectoryMode=0755
Environment="UROAM_LOG=info"

[Install]
WantedBy=multi-user.target